# _usr_lib_systemd_system_brew-packages.service
[Unit]
Description=Update Homebrew
After=network-online.target
Wants=network-online.target
After=brew-setup.service
ConditionPathExists=/var/home/linuxbrew/.linuxbrew/bin/brew

[Service]
Type=oneshot
#ExecStart=/usr/bin/bash -c "curl -fsSL https://raw.githubusercontent.com/Homebrew/install/HEAD/install.sh | bash"
ExecStart=/usr/bin/just --justfile /usr/share/just/justfile homebrew-packages

[Install]
WantedBy=multi-user.target
